# src/tools/monstache/vendor/github.com/coreos/go-systemd/fixtures/freeze.service
[Unit]
Description=freeze unit test

[Service]
ExecStart=/bin/sleep 400
